# 参考项目/declank/deploy/declank-tunnel.service
# systemd user service for the Cloudflare Tunnel that exposes the backend.
# Pairs with deploy/cloudflared.config.example.yml.
# Install:
#   cp deploy/declank-tunnel.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now declank-tunnel.service

[Unit]
Description=declank Cloudflare Tunnel
After=network-online.target declank.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/cloudflared tunnel --no-autoupdate --config /home/YOU/.cloudflared/declank.yml run
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
